# systemd/mcp-hub.socket
[Unit]
Description=CSA MCP Hub Socket

[Socket]
# Legacy `%t/csa` path remains compatible via migration-created symlink.
ListenStream=%t/cli-sub-agent/mcp-hub.sock
SocketMode=0600
DirectoryMode=0700
Service=mcp-hub.service

[Install]
WantedBy=sockets.target
